# debian/docprinter.service
[Unit]
Description=DocPrinter HTTP service (Docker)
Documentation=https://github.com/
After=docker.service containerd.service network-online.target
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=-/etc/default/docprinter
StandardOutput=journal
StandardError=journal
ExecStartPre=/usr/lib/docprinter/remove-container.sh
ExecStart=/usr/lib/docprinter/run-container.sh
ExecStopPost=-/usr/lib/docprinter/remove-container.sh
KillMode=mixed
TimeoutStopSec=60
Restart=on-failure
RestartPreventExitStatus=125
RestartSec=15

[Install]
WantedBy=multi-user.target
